# backend/config/crontab
# Run ingest job daily at 4:00 AM UTC
0 4 * * * cd /app && npm run ingest